# sw/slc7_x86-64/simulation/v1.0-1/etc/profile.d/init.sh
[ "X$DEFAULTS_RELEASE_VERSION" = X ] && source ${WORK_DIR}/slc7_x86-64/defaults-release/v1-1/etc/profile.d/init.sh
export SIMULATION_ROOT=/pool/condor/dir_55582/SHiPBuild/FairShip/sw/slc7_x86-64/simulation/v1.0-1
export SIMULATION_VERSION=v1.0
export SIMULATION_REVISION=1
export SIMULATION_HASH=07dfc6e33cac1c179c602988331328054030a13c
export SIMULATION_COMMIT=0
export DYLD_LIBRARY_PATH=$SIMULATION_ROOT/lib:$DYLD_LIBRARY_PATH
export PATH=$SIMULATION_ROOT/bin:$PATH
export LD_LIBRARY_PATH=$SIMULATION_ROOT/lib:$LD_LIBRARY_PATH
